# scripts/launchd/com.lmw.youtube-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lmw.youtube-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/ExternalLiumw/lavori/01_code/product-simulator-media-capture/scripts/youtube-monitor-daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Volumes/ExternalLiumw/lavori/01_code/product-simulator-media-capture</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/liumingwei</string>
        <key>HTTP_PROXY</key>
        <string>http://127.0.0.1:7897</string>
        <key>HTTPS_PROXY</key>
        <string>http://127.0.0.1:7897</string>
        <key>ALL_PROXY</key>
        <string>http://127.0.0.1:7897</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Volumes/ExternalLiumw/lavori/01_code/product-simulator-media-capture/data/youtube_monitor/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/ExternalLiumw/lavori/01_code/product-simulator-media-capture/data/youtube_monitor/logs/launchd-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
